# neptune/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:python]
command=python accessible_test/client.py %(ENV_HERCULES_CONTAINER_IP)s 80
